# stream-camera-listen.service
[Unit]
Description=Camera MJPEG: UDP discovery and TCP server until client connects
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pavel
WorkingDirectory=/home/pavel/projects/Mobile_Raspberry_5-
ExecStart=/home/pavel/projects/Mobile_Raspberry_5-/.venv/bin/python /home/pavel/projects/Mobile_Raspberry_5-/stream_camera.py --log-level INFO send --listen --port 5000
Restart=on-failure
RestartSec=4

[Install]
WantedBy=multi-user.target
